# airflow-image/supervisord.conf
[supervisord]
nodaemon=true

[program:scheduler]
command=airflow scheduler -D

[program:webserver]
command=airflow webserver -p 8080 -D